# main.cron
10 1 * * 1-5 /home/deaian/kshn/project/main.py